# etc/supervisor/conf.d/flask_app.conf
[program:flask_app]
directory=/app
command=gunicorn wsgi:application --workers=3
user=www-data
autostart=true
autorestart=true
environment=
    FLASK_APP=app.flask_app,
    FLASK_ENV=production