# etc/supervisor/conf.d/flask_app.conf
[program:app]
directory=/
command=gunicorn wsgi:application --workers=3
user=www-data
autostart=true
autorestart=true
environment=
    FLASK_APP=app,
    FLASK_ENV=production